# bsbf-cdc-ether/files/etc/hotplug.d/usb/99-bsbf-cdc-ether
#!/bin/sh

[ "$ACTION" = "bind" ] && [ "$DEVTYPE" = "usb_interface" ] && bsbf-cdc-ether
